# nemoclaw/nemoclaw-sentinel.service
# NemoClaw Sentinel — 自主巡檢 systemd unit(persistent deployment)
# 安裝:sudo cp nemoclaw/nemoclaw-sentinel.service /etc/systemd/system/
#       sudo systemctl daemon-reload && sudo systemctl enable --now nemoclaw-sentinel
# 檢視:journalctl -u nemoclaw-sentinel -f
[Unit]
Description=NemoClaw Sentinel autonomous multimodal safety patrol (GB10)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
# 將 User 改成你的帳號;%h 會解析為該帳號的家目錄(repo 預設置於 ~/Security-AI-Agent)
User=sentinel
WorkingDirectory=%h/Security-AI-Agent
ExecStart=/usr/bin/env bash %h/Security-AI-Agent/nemoclaw/nemoclaw-supervisor.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
